# litellm/.config/systemd/user/litellm-proxy.service
[Unit]
Description=LiteLLM Proxy (Tensorix cache layer)
After=network.target

[Service]
Type=simple
# Pull the Tensorix key from `pass` at start time. Secret is never written to disk.
ExecStart=/bin/sh -c 'TENSORIX_API_KEY=$(pass show api/tensorix) exec %h/.local/litellm-proxy/venv/bin/litellm --config %h/.config/litellm/config.yaml --port 4000 --host 127.0.0.1'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
